# collections/ansible_collections/netfusion/monitoring/roles/prometheus_pve_exporter/templates/prometheus-pve-exporter.service.j2
[Unit]
Description=Prometheus Proxmox VE Exporter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ pve_exporter_user }}
Group={{ pve_exporter_group }}
WorkingDirectory={{ pve_exporter_home }}
Environment=HOME={{ pve_exporter_home }}
ExecStart={{ pve_exporter_venv }}/bin/pve_exporter --config.file {{ pve_exporter_config_file }} --web.listen-address {{ pve_exporter_listen_address }}{% for flag in pve_exporter_collectors %} {{ flag }}{% endfor %}

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target

